# supervisord.conf
[supervisord]
nodaemon=true

[program:myapp]
command=python main.py
autostart=true
autorestart=true
